# deploy/dhcp-toggle-webui.service
[Unit]
Description=DHCP Toggle Web UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/usr/local/share/dhcp-toggle/webui
ExecStart=/usr/bin/python3 app.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
